# configs/systemd/timers/users/feh-bg@yukiteru.timer
[Unit]
Description=Change bg with feh hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
