# configs/systemd/services/system/volfread-hacks.service
[Unit]
Description=Volfread hacks service

[Service]
ExecStart=/usr/local/bin/volfread-hacks.sh

[Install]
WantedBy=default.target